# hadoop/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:hadoopdfs]
command=/bin/bash -c "start-dfs.sh"
